# sandbox/supervisord.conf
; Sandbox container supervisor — keeps the X stack, browser, and helper
; daemons alive so the Desktop / Browser tabs in TeamWork don't go dark
; when one of the underlying processes crashes.
;
; Process tree at runtime:
;   PID 1: supervisord
;     ├─ Xvfb (priority 10)
;     ├─ dbus session bus (priority 15)
;     ├─ x11vnc (priority 20) — needs Xvfb
;     ├─ xfwm4 / xfce4-panel / xfdesktop (priority 25-35) — need Xvfb + dbus
;     ├─ websockify (priority 40) — proxies noVNC :6080 → :5900
;     ├─ clipboard-bridge (priority 45)
;     ├─ chromium (priority 50) — CDP on :9222
;     ├─ socat (priority 55) — exposes :9222 as :9223 for cross-container access
;     └─ opencode (priority 60) — coding-agent server on :4096
;
; Logs are forwarded to docker logs via /dev/fd/1 and /dev/fd/2 so
; `docker compose logs sandbox` keeps showing what we had before.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
environment=DISPLAY=":99",DBUS_SESSION_BUS_ADDRESS="unix:path=/run/dbus-session.sock"

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── X stack ──────────────────────────────────────────────────────────

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1920x1080x24
priority=10
autostart=true
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:dbus]
command=/usr/bin/dbus-daemon --session --address=unix:path=/run/dbus-session.sock --nofork --syslog
priority=15
autostart=true
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -nopw -rfbport 5900 -q
priority=20
autostart=true
autorestart=true
startretries=10
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:xfwm4]
command=/usr/bin/xfwm4
priority=25
autostart=true
autorestart=true
startretries=10
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:xfce4-panel]
command=/usr/bin/xfce4-panel
priority=30
autostart=true
autorestart=true
startretries=10
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:xfdesktop]
command=/usr/bin/xfdesktop
priority=35
autostart=true
autorestart=true
startretries=10
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── Web bridges ──────────────────────────────────────────────────────

[program:websockify]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 localhost:5900
priority=40
autostart=true
autorestart=true
startretries=10
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:clipboard-bridge]
command=/usr/bin/python3 /usr/local/bin/clipboard-bridge.py
priority=45
autostart=true
autorestart=true
startretries=10
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── Browser ──────────────────────────────────────────────────────────

[program:chromium]
command=/usr/local/bin/chromium-launch.sh
priority=50
autostart=true
autorestart=true
startretries=10
startsecs=3
; Chrome forks zygote/gpu/renderer subprocesses; tear them all down on
; restart or supervisord shutdown so we don't pile up <defunct> entries.
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:cdp-proxy]
command=/usr/bin/socat TCP-LISTEN:9223,fork,reuseaddr,bind=0.0.0.0 TCP:127.0.0.1:9222
priority=55
autostart=true
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── Coding-agent server ──────────────────────────────────────────────

[program:opencode]
command=/usr/local/bin/opencode serve --hostname 0.0.0.0 --port 4096
priority=60
autostart=true
autorestart=true
startretries=10
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
